# 27-udevcrypt.rules.in
ACTION=="add", KERNEL =="sd*", TYPE=="disk", RUN+="@prefix@/bin/cryptcheck"
