# fsarchiver-backup.service.in
[Unit]
Description=Makes fsarchiver backup of entire system.
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=@bindir@/fsarchiver-backup

[Install]
WantedBy=shutdown.target reboot.target

