# pwd-force-locker.service.in
[Unit]
Description=Forcefully sets passwords expired.
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=@BINDIR@/pwd-force-locker

[Install]
WantedBy=shutdown.target reboot.target

